# gateway/etc/init.d/krill-tc
#!/sbin/openrc-run

description="Krill traffic control"

depend() {
    need net
    after nftables
}

start() {
    ebegin "Starting Krill traffic control"
    /etc/tc.qos
    eend $?
}

stop() {
    ebegin "Stopping Krill traffic control"
    # Clear the queues on both interfaces.
    
    tc qdisc del dev eth0 root 2>/dev/null || true
    tc qdisc del dev wlan0 root 2>/dev/null || true
    eend 0
}
